# linux/xremap.service
[Unit]
Description=xremap

[Service]
Restart=always
ExecStart=xremap /usr/local/etc/xremap.yaml --watch

[Install]
WantedBy=multi-user.target
