# linux/xremap.service
[Unit]
Description=xremap

[Service]
Restart=always
ExecStart=xremap /home/k/Desktop/a/xremap.yaml --watch

[Install]
WantedBy=multi-user.target
